# .github/service/telethon-service-monitor.service
[Unit]
Description=Telethon Service Monitor
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/telegram-listener
ExecStart=/usr/bin/python3 service_monitor.py
EnvironmentFile=/root/telegram-listener/.env
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/root/telegram-listener
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target